# security/heimdal/files/smf/manifest.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--

NOTE:  This service manifest is not editable; its contents will
be overwritten by package or patch operations, including
operating system upgrade.  Make customizations in a different
file.
-->

<service_bundle type='manifest' name='@SMF_NAME@'>

  <service
      name='@SMF_PREFIX@/@SMF_NAME@/kdc'
      type='service'
      version='1'>
    <create_default_instance  enabled="false" />
    <dependency
	name='milestone-network'
	grouping='require_all'
	restart_on='error'
	type='service'>
      <service_fmri value='svc:/milestone/network' />
    </dependency>

    <exec_method
	type='method'
	name='start'
	exec='@PREFIX@/libexec/kdc --detach'
	timeout_seconds='60' />

    <exec_method
	type='method'
	name='stop'
	exec=':kill'
	timeout_seconds='60' />

    <stability value='Evolving' />
    <template>
      <common_name>
        <loctext xml:lang="C">
          Heimdal Kerberos Key Distribution Center
          </loctext>
        </common_name>
      </template>
  </service>
  <service
      name='@SMF_PREFIX@/@SMF_NAME@/kadmin'
      type='service'
      version='1'>
    <create_default_instance  enabled="false" />
    <dependency
	name='milestone-network'
	grouping='require_all'
	restart_on='error'
	type='service'>
      <service_fmri value='svc:/milestone/network' />
    </dependency>

    <exec_method
	type='method'
	name='start'
	exec='@PREFIX@/libexec/kadmind &amp;'
	timeout_seconds='60' />

    <exec_method
	type='method'
	name='stop'
	exec=':kill'
	timeout_seconds='60' />

    <stability value='Evolving' />
    <template>
      <common_name>
        <loctext xml:lang="C">
          Heimdal Kerberos Administration Service
          </loctext>
        </common_name>
      </template>
  </service>

  <service
      name='@SMF_PREFIX@/@SMF_NAME@/kpasswd'
      type='service'
      version='1'>
    <create_default_instance  enabled="false" />
    <dependency
	name='milestone-network'
	grouping='require_all'
	restart_on='error'
	type='service'>
      <service_fmri value='svc:/milestone/network' />
    </dependency>

    <exec_method
	type='method'
	name='start'
	exec='@PREFIX@/libexec/kpasswdd &amp;'
	timeout_seconds='60' />

    <exec_method
	type='method'
	name='stop'
	exec=':kill'
	timeout_seconds='60' />

    <stability value='Evolving' />
    <template>
      <common_name>
        <loctext xml:lang="C">
          Heimdal Kerberos Password Change Daemon
          </loctext>
        </common_name>
      </template>
  </service>

  <service
      name='@SMF_PREFIX@/@SMF_NAME@/iprop'
      type='service'
      version='1'>
    <instance name="master" enabled="false">
      <dependency
	  name='milestone-network'
	  grouping='require_all'
	  restart_on='error'
	  type='service'>
	<service_fmri value='svc:/milestone/network' />
      </dependency>
      <dependency
	  name='filesystem-local'
	  grouping='require_all'
	  restart_on='error'
	  type='service'>
	<service_fmri value='svc:/system/filesystem/local' />
      </dependency>
      <dependency
	  name='config_data'
	  grouping='require_all'
	  restart_on='restart'
	  type='path'>
	<service_fmri value='file://localhost/var/heimdal/slaves' />
      </dependency>
      <exec_method
	  type='method'
	  name='start'
	  exec='@PREFIX@/libexec/ipropd-master &amp;'
	  timeout_seconds='60' />
      <exec_method
	  type='method'
	  name='stop'
	  exec=':kill'
	  timeout_seconds='60' />
    </instance>
    <instance name='slave' enabled='false'> 
      <dependency
	  name='milestone-network'
	  grouping='require_all'
	  restart_on='error'
	  type='service'>
	<service_fmri value='svc:/milestone/network' />
      </dependency>
      <dependency
	  name='filesystem-local'
	  grouping='require_all'
	  restart_on='error'
	  type='service'>
	<service_fmri value='svc:/system/filesystem/local' />
      </dependency>
      <exec_method
	  type='method'
	  name='start'
	  exec='@PREFIX@/libexec/ipropd-slave %{config/master-kdc} &amp;'
	  timeout_seconds='60' />
      <exec_method
	  type='method'
	  name='stop'
	  exec=':kill'
	  timeout_seconds='60' />
      <property_group name='config' type='application'>
	<propval name='master-kdc'
		 type='astring'
		 value='__REPLACE_ME__' />
      </property_group>
    </instance>
    <stability value='Evolving' />
    <template>
      <common_name>
	<loctext xml:lang="C">
	  Heimdal Incremental Propagation Service
	  </loctext>
	</common_name>
      </template>
  </service>
</service_bundle>
